# launchd/com.example.codex.deepseek-responses-adapter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.codex.deepseek-responses-adapter</string>
  <key>ProgramArguments</key>
  <array>
    <string>/absolute/path/to/node</string>
    <string>/absolute/path/to/codex-deepseek-responses-adapter/bin/codex-deepseek-responses-adapter.mjs</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/codex-deepseek-responses-adapter.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/codex-deepseek-responses-adapter.log</string>
</dict>
</plist>
